# home/.config/autostart/r3dsh1ft.desktop
[Desktop Entry]
Type=Application
Name=r3dsh1ft
Exec=redshift -P -O 2000
